# .config/autostart/xset.desktop
[Desktop Entry]
Type=Application
Name=Xset
Comment=Move faster
Exec=xset r rate 300 50
StartupNotify=false
