# .config/autostart/xset.desktop
[Desktop Entry]
Type=Application
Name=xset
Comment=Move faster
Exec=xset r rate 300 50
StartupNotify=false
Terminal=false
Hidden=false
